# coredns/root/etc/services.d/coredns/run
#!/command/with-contenv bash

/app/coredns -dns.port=53 -conf /app/config/coredns/Corefile